# service/tinstaller.service.example
[Unit]
Description=TInstaller - Unified APK Distribution Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/opt/web-serv
Environment="PATH=/opt/web-serv/venv/bin"
EnvironmentFile=/opt/web-serv/.env
ExecStart=/opt/web-serv/venv/bin/python3 /opt/web-serv/server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
